# yocto/meta-homeagent/recipes-connectivity/zigbee2mqtt/files/zigbee2mqtt.service
[Unit]
Description=Zigbee2MQTT Bridge
After=network.target mosquitto.service
Wants=mosquitto.service

[Service]
Type=simple
WorkingDirectory=/opt/zigbee2mqtt
EnvironmentFile=-/etc/default/zigbee2mqtt
ExecStart=/usr/bin/node /opt/zigbee2mqtt/index.js
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

# Security hardening
User=root
Group=root

[Install]
WantedBy=multi-user.target
